# setup/systemd/ehecatl.service
[Unit]
Description=Ehecatl OpenCore Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__/app
ExecStart=/usr/bin/env node index.js
Restart=always
RestartSec=2
TimeoutStopSec=30
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
